# ps3gamepad.service
 [Unit]  
Description=Sets up the PS3 gamepad.  
[Service]
Type=oneshot
ExecStart=/bin/bash /home/xfce/Documents/ps3/slaveit.sh